# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:flask]
command=./bin/manage.py runserver
redirect_stderr=true

[program:smtpd]
command=python -um smtpd -nc DebuggingServer
redirect_stderr=true

[program:celery]
command=python -u .env/bin/celery worker -A src.tasks.worker.celery -l INFO
redirect_stderr=true
